# etc/systemd/system/yay-reindex.service
[Unit]
Description=Reindex the yay package database

[Service]
Type=oneshot
ExecStart=yay -Pcc
